# osx/Library/LaunchAgents/me.kucho.environment.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>me.kucho.environment</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>
            launchctl setenv XDG_CONFIG_HOME /Users/kucho/.config &&
            launchctl setenv XDG_CACHE_HOME /Users/kucho/.cache &&
            launchctl setenv XDG_DATA_HOME /Users/kucho/.local/share &&
            launchctl setenv XDG_STATE_HOME /Users/kucho/.local/state &&
            launchctl setenv XDG_RUNTIME_DIR /Users/kucho/.local/run &&
            launchctl setenv XDG_BIN_HOME /Users/kucho/.local/bin
        </string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>ServiceIPC</key>
    <false/>
</dict>
</plist>
